# roles/graphite/files/supervisord/graphite-api.conf
[program:graphite-api]
command = gunicorn -w2 graphite_api.app:app -b 0.0.0.0:8888
autostart = true
autorestart = true
stopsignal = kill