# roles/hoardd/files/supervisor/hoardd.conf
[program:hoardd]
command=node /opt/hoardd/start.js 
process_name=%(program_name)s
directory= /opt/hoardd/
autostart=true
autorestart=true
stopsignal=QUIT
