# roles/statsd/files/supervisor/statsd.conf
[program:statsd]
command=node /opt/statsd/stats.js /opt/statsd/local_config.js
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT